# dist/debian/obdi-worker.upstart
# obdi - script runner worker

description     "script runner worker"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec obdi-worker

